# unifi-protect-backup/rootfs/etc/services.d/unifi-protect-backup/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting UniFi Protect Backup"

ip_address=$(bashio::config "ip_address")
username=$(bashio::config "username")
password=$(bashio::config "password")
retention=$(bashio::config "retention")
rclone_destination=$(bashio::config "rclone_destination")
rclone_configuration=$(bashio::config "rclone_configuration")

rclone_config_file_path=$(rclone config file | cut -d ":" -f 2 | xargs)
echo "$rclone_configuration" > "$rclone_config_file_path"

verify_ssl=$(bashio::config "verify_ssl")
if [ "$verify_ssl" = "false" ]; then
    verify_ssl="--no-verify-ssl"
else
    verify_ssl="--verify-ssl"
fi

skip_missing=$(bashio::config "skip_missing")
if [ "$skip_missing" = "true" ]; then
    skip_missing="--skip-missing"
else
    skip_missing=""
fi

extra_args=$(bashio::config "extra_args")

unifi-protect-backup \
    --address "$ip_address" \
    --username "$username" \
    --password "$password" \
    --retention "$retention" \
    --rclone-destination "$rclone_destination" \
    --sqlite_path "/data/unifi_protect_backup.sqlite" \
    $verify_ssl \
    $skip_missing \
    $extra_args
